# systemd/website-cloner.service.template
[Unit]
Description=Website Cloner - Static Site Cloning Tool
After=network.target oauth2-proxy.service
Wants=oauth2-proxy.service

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory={{WORKING_DIRECTORY}}
ExecStart=/usr/bin/node {{WORKING_DIRECTORY}}/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=website-cloner

# Environment
Environment=NODE_ENV=production
Environment=PORT=3000

# Security settings
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
